# mastodon/mastodon-database-backup.timer
[Unit]
Description=Backup mastodon database with restic daily
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=multi-user.target
