# Pixelfed/pixelfed-backup.timer
[Unit]
Description=Backup pixelfed with restic daily
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=multi-user.target
